# voice-agent/runpod/supervisord.conf
; supervisord config for the all-in-one RunPod image.
; Logs all go to stdout/stderr so RunPod's log viewer captures everything.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

; ─── Ollama LLM daemon ────────────────────────────────────────────────────────
[program:ollama]
command=/usr/local/bin/ollama serve
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=OLLAMA_MODELS="/workspace/ollama-models",OLLAMA_HOST="127.0.0.1:11434"

; ─── Kokoro TTS server ────────────────────────────────────────────────────────
; The base kokoro-fastapi-gpu image installs the server at /app and exposes uvicorn
; via "api.src.main:app" on port 8880. If a future image rev changes that path, run
; `docker run --rm --entrypoint cat ghcr.io/remsky/kokoro-fastapi-gpu:latest /app/Dockerfile`
; to find the new CMD and update this line.
[program:kokoro]
command=python -m uvicorn api.src.main:app --host 0.0.0.0 --port 8880
directory=/app
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ─── Pipecat agent ────────────────────────────────────────────────────────────
; startsecs=15 lets Ollama + Kokoro warm up before we mark the agent "started".
; If the agent dies during that window supervisord considers it a failed start
; and logs the traceback prominently.
[program:agent]
command=python -m agent.main
directory=/opt/agent
autostart=true
autorestart=true
startretries=10
startsecs=15
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
